# apps/onboarding-system/deploy/onboarding-system.service
[Unit]
Description=Onboarding System
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=hnester
Group=hnester
WorkingDirectory=/home/hnester/onboarding-system
Environment=HOME=/home/hnester
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=NODE_ENV=development
ExecStart=/usr/bin/npm run dev
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5
LimitNOFILE=4096
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
